# deploy/llamesh.service
[Unit]
Description=llamesh
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=%REPO_PATH%
ExecStart=%REPO_PATH%/target/release/llamesh
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=3600
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llamesh

[Install]
WantedBy=default.target
